# Library/LaunchAgents/pl.com.siema.env.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>Setup ENV vars</string>
	<key>ProgramArguments</key>
	<array>
		<string>sh</string>
		<string>-c</string>
		<string>
            launchctl setenv PATH /usr/local/bin:$(launchctl getenv PATH)
            launchctl setenv ANDROID_HOME $HOME/android/sdk
            launchctl setenv NDK_HOME $HOME/android/sdk/ndk-bundle
            launchctl setenv AW_ENGINE_PATH $HOME/projects/multi/AWEngine
            launchctl setenv AW_PLUGIN_MANAGER_PATH $HOME/projects/multi/AWPluginManager
            launchctl setenv TC_RESKIN_PATH $HOME/projects/multi/TCReskin
        </string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
